# board/diasom/overlay/etc/profile.d/local.sh
export HISTFILE='/dev/null'
export INPUTRC='/etc/inputrc'
export LANG='en_US'
export LC_ALL='C'
export PATH=$PATH':/usr/local/bin:/usr/local/sbin'
export TERM='vt102'
export TMPDIR='/var/tmp'
export TZ='UTC'

# Unhide console cursor
echo -e -n "\e[?25h"
